# elastic-kibana/services/kibana/run
#!/bin/sh
exec su - elastic -c "/opt/kibana/bin/kibana"
