# elastic-kibana/services/elastic/run
#!/bin/sh
exec su - elastic -c "/opt/elasticsearch/bin/elasticsearch -Edefault.transport.host=127.0.0.1 -Edefault.http.host=0.0.0.0"
